# src/github.com/bajratech/webhook/webhook.conf
[program:goserver]
command=/usr/bin/webhook.sh
autostart=true
autorestart=true