# src/github.com/bajratech/webhook/webhook.conf
[program:goserver]
command=/usr/bin/webhook.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/webhook.err.log
stdout_logfile=/var/log/webhook.out.log